# borgmatic.d/borgmatic.timer
[Unit]
Description=Run borgmatic backup

[Timer]
# Daily at 3am
OnCalendar=*-*-* 03:00
Persistent=true

[Install]
WantedBy=timers.target

